# systemd/ddns.service
[Unit]
Description=Update dynamic DNS
After=network.target network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/ddns
User=nobody
